# config/com.accuport.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.accuport.dashboard</string>
    <key>ProgramArguments</key>
    <array>
        <string>INSTALL_PATH/venv/bin/python</string>
        <string>INSTALL_PATH/dashbored/app.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>INSTALL_PATH/dashbored</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>INSTALL_PATH/dashbored/app.log</string>
    <key>StandardErrorPath</key>
    <string>INSTALL_PATH/dashbored/app.log</string>
</dict>
</plist>
